# supervisor-nuc1.conf
[program:nuc1]
command=/usr/local/bin/nuc1
stdout_logfile=/var/log/nuc1.log
